# debian/prerm
#!/bin/sh
if [ "$1" = "purge" ]; then
	rm -f /etc/twisted-web/disable
fi
invoke-rc.d twisted-web stop
invoke-rc.d twisted-web-users stop
